# root/etc/services.d/optimiser/run
#!/usr/bin/with-contenv bash

export HOME=/config

exec s6-setuidgid docker /usr/bin/python3 /optimiser
